# etc/inittab
# /etc/inittab

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::sysinit:swapon -a

tty1::respawn:su - juo

::ctrlaltdel:reboot

::shutdown:kill -15 -1
::shutdown:echo '*** SHUTING DOWN ***'
::shutdown:sleep 10
::shutdown:kill -9  -1
::shutdown:swapoff -a
::shutdown:umount -a -r

::once:ntpd -dd -q -p pool.ntp.org
::once:chown juo /sys/class/backlight/*/brightness
::once:chown juo /dev/tty1

::respawn:wpa_supplicant -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf
